# etc/udev/rules.d/70-persistent-monitor.rules
KERNEL=="card0",SUBSYSTEM=="drm",ACTION=="change",RUN+="/bin/su - -c /usr/local/bin/configure-monitors brayan"
